# intranet/backup/crontab
# m h  dom mon dow   command
0 0 * * 6 /root/scripts/full_bkp.sh
55 23 * * 0,1,2,3,4,5 /root/scripts/incr_bkp.sh
